# crontab
*/10 * * * * /app/run-tests.sh >> /app/cron.log 2>&1
